# config/systemd.service
[Unit]
Description=Python MTProxy Service
Documentation=https://github.com/OrangeKing92/mtproxy-manager
After=network.target network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=mtproxy
Group=mtproxy
WorkingDirectory=/opt/python-mtproxy

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/python-mtproxy
Environment=PYTHONUNBUFFERED=1

# Main process
ExecStart=/usr/bin/python3 -m mtproxy.server
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security settings
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true

# File system access
ReadWritePaths=/opt/python-mtproxy
ReadOnlyPaths=/etc/ssl/certs

# Network access
PrivateNetwork=false

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitCORE=0

# Memory and CPU
MemoryAccounting=true
CPUAccounting=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mtproxy

[Install]
WantedBy=multi-user.target
